# .crontab
# m h  dom mon dow   command
* * * * * /home/g/MyScripts/autoECN/ECNlogger /home/g/MyScripts/autoECN
